# worker/systemd/frp-generate.timer
[Unit]
Description=Run frp-generate every 3 minutes

[Timer]
# Wait 2 min after boot so the system settles and ollama warms.
OnBootSec=2min
# 3 min between runs — gives an old 8 GB laptop breathing room and keeps
# ollama's model unload cycle (OLLAMA_KEEP_ALIVE=30s) working. Bump to 60s
# once you're confident the box handles it.
OnUnitActiveSec=3min
AccuracySec=15s
Unit=frp-generate.service

[Install]
WantedBy=timers.target
